# infra/vm/systemd/rucksack-autopilot-erniesg-paillette-drain.timer
# Generated by rucksack. Runs in the trusted VM user's local timezone.
[Unit]
Description=Nightly Rucksack autopilot drain
Documentation=man:systemd.timer(5)

[Timer]
OnCalendar=*-*-* 02:17:00
Persistent=true
RandomizedDelaySec=15m
Unit=rucksack-autopilot-erniesg-paillette-drain.service

[Install]
WantedBy=timers.target
